# scripts/watchparty-password-sync.service
[Unit]
Description=Watch Party RDS Password Sync
Documentation=https://github.com/EL-HOUSS-BRAHIM/watch-party
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/watch-party
ExecStart=/home/deploy/watch-party/scripts/sync-rds-password.sh
StandardOutput=journal
StandardError=journal

# Environment variables (if needed)
Environment=AWS_DEFAULT_REGION=eu-west-3
Environment=AWS_REGION=eu-west-3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/deploy/watch-party

[Install]
WantedBy=multi-user.target